# systemd/boss-zhipin.service.template
[Unit]
Description=BOSS Zhipin daily job collection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=__WORK_DIR__
ExecStart=__WORK_DIR__/run_daily.sh
